# solar_optimiser/rootfs/etc/services.d/solar-optimiser/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the solar optimiser service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Starting solar optimiser service..."

# Check if Python is available
if ! command -v python3 &> /dev/null; then
    bashio::log.error "Python3 is not installed!"
    exit 1
fi

# Check if the program exists
if [ ! -f /usr/bin/my_program ]; then
    bashio::log.error "my_program not found at /usr/bin/my_program"
    exit 1
fi

# Check if the program is executable
if [ ! -x /usr/bin/my_program ]; then
    bashio::log.error "my_program is not executable"
    exit 1
fi

# Check if templates directory exists
if [ ! -d /usr/bin/templates ]; then
    bashio::log.error "templates directory not found at /usr/bin/templates"
    exit 1
fi

# Export required environment variables
export SUPERVISOR_TOKEN="${SUPERVISOR_TOKEN}"
export HASS_URL="${HASS_URL}"

bashio::log.info "Starting Python web server..."
exec python3 /usr/bin/my_program
